# python-image/cron.d/cron_python
HOSTNAME=1bd340ee8998
PYTHON_VERSION=3.8.0
PWD=/scripts
TZ=Asia/Tokyo
HOME=/root
LANG=C.UTF-8
GPG_KEY=E3FF2839C048B25C084DEBE9B26995E310250568
TERM=xterm
SHLVL=1
PYTHON_PIP_VERSION=19.3.1
PYTHON_GET_PIP_SHA256=b86f36cc4345ae87bfd4f10ef6b2dbfa7a872fbff70608a1e43944d283fd0eee
PYTHON_GET_PIP_URL=https://github.com/pypa/get-pip/raw/ffe826207a010164265d9cc807978e3604d18ca0/get-pip.py
PATH=/usr/local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

*/1 *   * * *   root    mkdir /scripts/test
*/1 *   * * *   root    /usr/local/bin/python /scripts/file_make.py
